# cloud-api/deploy/sanfeng-cloud-api-backup.timer
[Unit]
Description=Daily Sanfeng Finance Cloud API data backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
Unit=sanfeng-cloud-api-backup.service

[Install]
WantedBy=timers.target
